# initrd/etc/udev/rules.d/90-fstab-automount.rules
ACTION=="add|change|remove",SUBSYSTEM=="block",KERNEL=="sd*" \
RUN+="/bin/sh -c '/sbin/make-fstab --update'"
